# etc/periodic/hourly/hourly
#!/bin/sh

current_time=$(date +"%H:%M:%S")
echo "Après 1 heure : L'heure est : $current_time"