# etc/periodic/monthly/mounthly
#!/bin/sh

current_time=$(date +"%H:%M:%S")
echo "Après 1 mois, L'heure actuelle est : $current_time"