# etc/periodic/15min/15min
#!/bin/sh

current_time=$(date +"%H:%M:%S")
echo "Après 15 min : L'heure est : $current_time"